# systemd/blinkenshell-sync-clone.service
[Unit]
Description=Sync the clone blinkenshell website from the main server
RequiresMountsFor=/mnt/efs
 
[Service]
Type=simple
ExecStart=/var/www/aws/systemd/blinkenshell-sync-clone
Restart=no
User=ubuntu
Group=ubuntu
