# systemd/git-sync-clone.timer
[Unit]
Description=Sync the clone git repos timer

[Timer]
OnCalendar=*-*-* 23:58:00

Persistent=true

[Install]
WantedBy=timers.target
