# systemd/Zz-sync-clone.timer
[Unit]
Description=Sync the clone Zz directory timer

[Timer]
OnCalendar=*-*-* 23:55:00

Persistent=true

[Install]
WantedBy=timers.target
