# etc/profile.d/ip-color.sh
# /etc/profile.d/ip-color.sh
alias ip='ip --color=auto'
